# sys-cluster/nova/files/systemd/nova-objectstore.service
[Unit]
Description=OpenStack Nova Object Store

[Service]
Type=simple
ExecStart=/usr/bin/nova-objectstore --flagfile=/etc/nova/nova.conf
StandardOutput=syslog
StandardError=syslog
Restart=always
User=nova

[Install]
WantedBy=multi-user.target
