# sys-cluster/nova/files/systemd/nova-volume.service
[Unit]
Description=OpenStack Nova Volume Service

[Service]
Type=simple
ExecStart=/usr/bin/nova-volume --flagfile=/etc/nova/nova.conf
StandardOutput=syslog
StandardError=syslog
Restart=always
User=nova

[Install]
WantedBy=multi-user.target
