# sys-cluster/nova/files/systemd/nova-cert.service
[Unit]
Description=OpenStack Nova Certificate Service

[Service]
Type=simple
ExecStart=/usr/bin/nova-cert --flagfile=/etc/nova/nova.conf
StandardOutput=syslog
StandardError=syslog
Restart=always
User=nova

[Install]
WantedBy=multi-user.target
